# docker/app/supervisor.conf
[supervisord]
nodaemon=true

[program:app-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
user=root
stdout_logfile=/var/www/app/storage/logs/worker.log

[program:segmentation_socket_server]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23456
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server.log

[program:ocsr_socket_server]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65432
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server.log

[program:classifier_socket_server]
command=python3 /var/www/app/app/Python/decimer_classifier_server.py 11111
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/classifier_socket_server.log


[program:segmentation_socket_server_23457]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_segmentation_server.py 23457"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23457.log

[program:segmentation_socket_server_23458]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_segmentation_server.py 23458"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23458.log

[program:segmentation_socket_server_23459]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_segmentation_server.py 23459"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23459.log

[program:segmentation_socket_server_23460]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_segmentation_server.py 23460"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23460.log

[program:segmentation_socket_server_23461]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_segmentation_server.py 23461"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23461.log

[program:ocsr_socket_server_65433]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_predictor_server.py 65433"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65433.log

[program:ocsr_socket_server_65434]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_predictor_server.py 65434"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65434.log

[program:decimer_classifier_socket_11112]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_classifier_server.py 11112"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/decimer_classifier_server_11112.log

[program:decimer_classifier_socket_server_11113]
command=/bin/bash -c "sleep 120 && python3 /var/www/app/app/Python/decimer_classifier_server.py 11113"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/decimer_classifier_server_11113.log


[program:clear_config]
process_name=clear_config
command=php /var/www/app/artisan config:clear
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/clear_config.log

[program:clear_cache]
process_name=clear_cache
command=php /var/www/app/artisan cache:clear
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/clear_cache.log

[program:laravel_link]
process_name=create_storage_link
command=php /var/www/app/artisan storage:link
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/storage_link.log